# src/menubar/com.brunowinter.monitor_cc_menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.brunowinter.monitor_cc_menubar</string>

  <!-- Replaced at setup time: points to Monitor_CC_Menubar.app/Contents/MacOS/Monitor_CC_Menubar (py2app native binary) -->
  <key>ProgramArguments</key>
  <array>
    <string><BUNDLE_LAUNCHER></string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd default PATH is /usr/bin:/bin:/usr/sbin:/sbin — no Homebrew -->
    <!-- arm64 Homebrew: /opt/homebrew/bin; Intel fallback: /usr/local/bin  -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/monitor_cc_menubar.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/monitor_cc_menubar.err</string>
</dict>
</plist>
