# src/menubar/com.brunowinter.monitor-cc-menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.brunowinter.monitor-cc-menubar</string>

  <!-- Replaced at setup time: points to monitor-cc-menubar.app/Contents/MacOS/monitor-cc-menubar (py2app native binary) -->
  <key>ProgramArguments</key>
  <array>
    <string><BUNDLE_LAUNCHER></string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd default PATH is /usr/bin:/bin:/usr/sbin:/sbin — no Homebrew -->
    <!-- arm64 Homebrew: /opt/homebrew/bin; Intel fallback: /usr/local/bin  -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- launchd sets no locale → Python defaults to ascii → UnicodeDecodeError on non-ASCII
         process names (e.g. worker spawn-prompts with emoji/umlauts in ps -A -o command= output).
         PYTHONUTF8=1 forces UTF-8 for all text=True subprocess calls + open() in this process.
         Belt-and-suspenders: per-call encoding='utf-8',errors='replace' is the primary fix. -->
    <key>PYTHONUTF8</key>
    <string>1</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/monitor-cc-menubar.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/monitor-cc-menubar.err</string>
</dict>
</plist>
